# Chapter_09/security-brief.timer
[Unit]
Description=Run CAH2 security brief each morning

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
